# airootfs/root/.config/autostart/nwg-dock.desktop
[Desktop Entry]
Name=Gory OS Dock
Comment=Wayland compatible dock
Exec=nwg-dock -d
Terminal=false
Type=Application
Categories=Utility;
NoDisplay=true
